# Scripts/local.tokenout.fetch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.tokenout.fetch</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Applications/Codex.app/Contents/Resources/cua_node/bin/node</string>
    <string>__HOME__/Library/Application Support/TokenOut/fetch-quota.js</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>60</integer>
  <key>StandardOutPath</key>
  <string>__HOME__/Library/Logs/TokenOut.fetch.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/Library/Logs/TokenOut.fetch.err.log</string>
</dict>
</plist>
